# maintainer-scripts/crontab
16  0 * * * sh /home/gccadmin/scripts/update_version_svn
50  0 * * * sh /home/gccadmin/scripts/update_web_docs_svn
55  0 * * * sh /home/gccadmin/scripts/update_web_docs_libstdcxx_svn
32 22 * * 3 sh /home/gccadmin/scripts/gcc_release -s 4.2:branches/gcc-4_2-branch -l -d /sourceware/snapshot-tmp/gcc all
32 22 * * 4 sh /home/gccadmin/scripts/gcc_release -s 4.3:branches/gcc-4_3-branch -l -d /sourceware/snapshot-tmp/gcc all
32 22 * * 5 sh /home/gccadmin/scripts/gcc_release -s 4.4:trunk           -l -d /sourceware/snapshot-tmp/gcc all
